# backend/metzscore-api.socket
[Unit]
Description=MetzScore FastAPI backend socket

[Socket]
ListenStream=127.0.0.1:8001
Accept=no

[Install]
WantedBy=sockets.target
